# deploy/roles/daoru/templates/supervisor.conf
[program:{{ daoru_name }}]
command=/usr/bin/php /work/Explorer/{{ daoru_name }}/main.php
autostart=true
autorestart=false
startretries=0
